# deploy/ansible/roles/services/templates/mbtawesome.conf
description "daemon for {{project_name}}"

start on startup
stop on shutdown
respawn

env NODE_ENV={{env}}

script
  chdir {{site_path}}
  /usr/bin/npm start
end script
